# etc/supervisord.conf
[unix_http_server]
file=%(here)s/../var/supervisor.sock

[supervisord]
logfile=%(here)s/../var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=20 
loglevel=info
pidfile=%(here)s/../var/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/../var/supervisor.sock

[program:celery]
command = %(here)s/../manage.py celeryd -l INFO -B
redirect_stderr = true
stdout_logfile = %(here)s/../var/log/%(program_name)s.log
stdout_capture_maxbytes=10MB
stdout_logfile_backups=20

[program:froide]
command = %(here)s/../manage.py runserver
redirect_stderr = true
stdout_logfile = %(here)s/../var/log/%(program_name)s.log
stdout_capture_maxbytes=10MB
stdout_logfile_backups=20
